# packaging/ubuntu/trusty/scripts/postrm.sh
#!/bin/sh
update-rc.d -f goiardi remove
